# package/base-files/files-nommu/etc/inittab
::sysinit:/etc/init.d/rcS S boot
::shutdown:/etc/init.d/rcS K shutdown
tts/0::askfirst:/bin/hush
ttyS0::askfirst:/bin/hush
tty1::askfirst:/bin/hush
